# templates/upstart.conf.j2
#!upstart

# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

description "Start celery flower for {{flower_app_name}}"

start on (local-filesystems and runlevel [2345])
stop on runlevel [!2345]

{% if flower_user == 'root' %}
env C_FORCE_ROOT=true
env HOME=/root
{% endif %}

env USER="{{ flower_user }}"
env GROUP="{{ flower_group }}"

env FLOWER_NAME=""
env FLOWER_APP="{{ flower_app_module }}"
env FLOWER_CHDIR="{{ flower_work_dir }}"
env FLOWER_LOG_FILE="{{ flower_log_dir }}/flower.log"
env FLOWER_LOG_LEVEL="{{ flower_log_level }}"
env FLOWER_PID_FILE="{{ flower_run_dir }}/flower.pid"
env FLOWER_BASIC_AUTH="{{ flower_basic_auth }}"
env FLOWER_URL_PREFIX="{{ flower_url_prefix }}"
env FLOWER_BROKER="{{ flower_broker }}"
env FLOWER_HOST="{{ flower_host }}"
env FLOWER_PORT="{{ flower_port }}"

{% set env = flower_env %}
{% for key in env %}
env {{key}}={{env[key]}}
{% endfor %}

console log
kill timeout 20
limit nofile 65536 65536
respawn

script

    {# Run any pre-exec commands we have defined, e.g `. /path/to/your/env_vars/.env`#}
    {% for line in flower_pre_exec %}
    {{line}}
    {% endfor %}

    exec flower                --pidfile="$FLOWER_PID_FILE" \
                               --logfile="$FLOWER_LOG_FILE" \
                               --loglevel="$FLOWER_LOG_LEVEL" \
                               --app="$FLOWER_APP" \
                               --workdir="$FLOWER_CHDIR" \
                               --uid="$USER" \
                               --gid="$GROUP" \
                               --broker="$FLOWER_BROKER" \
                               --address="$FLOWER_ADDRESS" \
                               --port="$FLOWER_PORT" \
                               {%  if flower_url_prefix_enabled %}
                               --url_prefix="$FLOWER_URL_PREFIX" \
                               {%  endif %}
                               {%  if flower_basic_auth_enabled %}
                               --basic_auth="$FLOWER_BASIC_AUTH"
                               {%  endif %}
end script